# deploy/holiday-picker-update.service
# Self-update for the Wheel of Wander holiday picker.
#
# Started by holiday-picker-update.path whenever the app's admin presses
# "Update & restart server" (the app writes an update-requested flag file
# into its state directory — it is sandboxed and can't run git/systemctl
# itself). Runs as root: pulls the latest version and restarts the app.
#
# Install together with the .path unit:
#   sudo cp deploy/holiday-picker-update.{service,path} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now holiday-picker-update.path

[Unit]
Description=Wheel of Wander — pull latest version and restart

[Service]
Type=oneshot
# Remove the flag first so a failed pull doesn't retrigger in a loop
ExecStartPre=/bin/rm -f /var/lib/holiday-picker/update-requested
ExecStart=/usr/bin/git -C /opt/holiday-picker pull --ff-only
ExecStart=/bin/systemctl restart holiday-picker
