# deploy/holiday-picker-update.path
# Watches for the update-requested flag the app's admin button
# drops into the state directory, then starts holiday-picker-update.service.
# See holiday-picker-update.service for install instructions.

[Unit]
Description=Wheel of Choice — watch for update requests

[Path]
PathExists=/var/lib/holiday-picker/update-requested

[Install]
WantedBy=multi-user.target
